# inventory-frontend.service
[Unit]
Description=Inventory Tracker Frontend Build
After=network.target

[Service]
Type=oneshot
User=inventory
Group=inventory
WorkingDirectory=/opt/inventory/frontend
ExecStart=/usr/bin/npx vite build
Environment=NODE_ENV=production
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
