# inventory-backend.service
[Unit]
Description=Inventory Tracker Backend API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/inventory/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
